# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true

[program:hornet]
command=python3 /var/www/html/hornet.py
redirect_stderr=true
